# slippers-optimized.service
[Unit]
Description=Slippers FastAPI - High Performance
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Slippers

# Environment
Environment="PATH=/home/ubuntu/Slippers/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONOPTIMIZE=2"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="PYTHONUNBUFFERED=1"

# Load environment variables from .env
EnvironmentFile=/home/ubuntu/Slippers/.env

# Performance-optimized command
ExecStart=/home/ubuntu/Slippers/venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --loop uvloop \
    --http httptools \
    --backlog 4096 \
    --limit-concurrency 10000 \
    --limit-max-requests 100000 \
    --timeout-keep-alive 5 \
    --no-access-log \
    --log-level warning

# Restart policy
Restart=always
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=30

# Resource limits (prevent memory leaks)
MemoryMax=2G
MemoryHigh=1.5G
TasksMax=4096

# Process limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/Slippers/app/static
ReadWritePaths=/home/ubuntu/Slippers/slippers.db

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slippers-api

[Install]
WantedBy=multi-user.target
